# systemd/autowit-api.service
[Unit]
Description=Autowit AI Platform — NestJS API
After=network.target docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
User=autowit
WorkingDirectory=/autowit/backend
EnvironmentFile=/autowit/.env
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autowit-api

[Install]
WantedBy=multi-user.target
